# playbooks/roles/bifrost-create-vm-nodes/templates/vbmcd.service.j2
[Unit]
Description=Virtual BMC daemon

[Service]
Type=simple
Restart=on-failure
ExecStart={{ bifrost_venv_dir }}/bin/vbmcd --foreground
